# deploy/jmailcode.sh
#chkconfig: -85 15
nohup java -jar mailcode.jar &
